# rootfs/etc/crontabs/git
# min   hour    day     month   weekday command
*/15    *       *       *       *       /usr/bin/git-sync
